# wakepc.service
#
# sudo cp wakepc.service /etc/systemd/system/wakepc.service
# sudo systemctl daemon-reload
# sudo systemctl enable --now wakepc.service
# systemctl status wakepc.service
# journalctl -u wakepc.service -f
#

[Unit]
Description=Streamlit WakePC app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alvise
Group=alvise
WorkingDirectory=/home/alvise/wakepc
Environment=HOME=/home/alvise
Environment=PATH=/home/alvise/miniconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/alvise/miniconda3/bin/streamlit run /home/alvise/wakepc/app.py --server.headless=true --server.address=0.0.0.0 --server.port=9000 --browser.gatherUsageStats=false
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
